# rpistepper.service
[Unit]
Description=RPIStepper - Raspberry Pi Stepper Motor Controller
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/RPIStepper
Environment=RADAR_DASHBOARD_URL=http://127.0.0.1:5060
Environment=OPTA_DASHBOARD_URL=http://127.0.0.1:5070
ExecStart=/usr/bin/python3 /home/pi/RPIStepper/dashboard_server.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
